# ToBkup-Dec2024/archive_name/data/home/mirfaisal/CBDL-VQ/Test_App_Jerman_Repo/Elina_Utilities_Harman_repo_05112024/Elina_Utilities/ethernet/IDCEvo-network/switch_config.service
[Unit]
Description=Switch configuration
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/switch_config_util

[Install]
WantedBy=multi-user.target

